# rustclaw.service
[Unit]
Description=RustClaw Service
Wants=network-online.target
After=network-online.target

[Service]
Type=forking
User=guagua
Group=guagua
WorkingDirectory=/home/guagua/rustclaw
Environment=HOME=/home/guagua
Environment=RUSTCLAW_RUNTIME_ENV_SCRIPT=/home/guagua/runtime_env_filled.sh
Environment=RUSTCLAW_CONFIG_PATH=/home/guagua/rustclaw/configs/config.toml
Environment=RUSTCLAW_MODEL_SELECT=0
Environment=RUSTCLAW_LOG_COLOR=0
Environment=RUST_LOG=info
ExecStart=/bin/bash /home/guagua/rustclaw/start-all-bin.sh release
ExecStop=/bin/bash /home/guagua/rustclaw/stop-rustclaw.sh
PIDFile=/home/guagua/rustclaw/.pids/clawd.pid
Restart=on-failure
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=45
KillMode=control-group
TasksMax=512

[Install]
WantedBy=multi-user.target
